# devel/drone-runner-ssh/pkg-plist
bin/drone-runner-ssh
etc/rc.d/drone-runner-ssh
