# devel/drone-runner-exec/pkg-plist
bin/drone-runner-exec
etc/rc.d/drone-runner-exec
