# devel/drone-runner-exec/files/drone-runner-exec.in
#!/bin/sh

# PROVIDE: drone_runner_exec
# REQUIRE: NETWORKING SYSLOG
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable drone_runner_exec:
#
#drone_runner_exec_enable="YES"

. /etc/rc.subr

name="drone_runner_exec"

rcvar="drone_runner_exec_enable"

load_rc_config $name

: ${drone_runner_exec_user:="nobody"}
: ${drone_runner_exec_enable:="NO"}
: ${drone_runner_exec_facility:="daemon"}
: ${drone_runner_exec_priority:="debug"}
: ${drone_runner_exec_rpc_proto:="http"}
: ${drone_runner_exec_rpc_host:="localhost"}
: ${drone_runner_exec_rpc_secret:="super-duper-secret"}

command="%%PREFIX%%/bin/drone-runner-exec"
procname="%%PREFIX%%/bin/drone-runner-exec"

pidfile="/var/run/${name}.pid"

start_cmd="${name}_start"

drone_runner_exec_start() {
	if [ "${drone_runner_exec_user}" = "nobody" ]; then
		echo "rc_var \"drone_runner_exec_user\" need to be set to a valid user"
		exit 1
	fi
	/usr/sbin/daemon %%DAEMONARGS%% \
		-u ${drone_runner_exec_user} -p ${pidfile} \
		/usr/bin/env -i \
		"CLIENT_DRONE_RPC_PROTO=${drone_runner_exec_rpc_proto}" \
		"CLIENT_DRONE_RPC_HOST=${drone_runner_exec_rpc_host}" \
		"CLIENT_DRONE_RPC_SECRET=${drone_runner_exec_rpc_secret}" \
		$command
}

run_rc_command "$1"
